# files/service/hookyd/run
#!/bin/sh -e
# exec 2>&1
export PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/gonano/sbin:/opt/gonano/bin:/data/sbin:/data/bin"
exec /opt/gonano/hookyd/hookyd.lua 2>&1
