# files/service/sshd/run
#!/bin/sh -e
exec 2>&1
exec /opt/gonano/sbin/sshd -eD -f /opt/gonano/etc/ssh/sshd_config
